# ledfx/rootfs/etc/services.d/ledfx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: ledfx
# Runs ledfx
# ==============================================================================

bashio::log.info "Starting ledfx..."

exec ledfx --offline --config /data/ledfx
