# ledfx/rootfs/etc/s6-overlay/s6-rc.d/ledfx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: ledfx
# Runs ledfx
# ==============================================================================

bashio::log.info "Starting ledfx..."

exec ledfx --offline --config /data/ledfx
